# testing/munge/munge.initd
#!/sbin/openrc-run

name="munge"
pidfile="/run/munged.pid"
command="/usr/sbin/munged"
command_args="$DAEMON_ARGS"
command_user="$USER:$GROUP"

depend() {
	after net clock
}

start_pre() {
	checkpath -q -d -m 0755 -o "$USER:$GROUP" /run/"$RC_SVCNAME"
}

start_post() {
	mv /run/"$RC_SVCNAME"/munged.pid "$pidfile"
	checkpath -q -f -m 0644 -o root:root "$pidfile"
}
